# supervisord/eventwitter.conf
[program:eventwitter]
user=jimshoe
directory=/home/jimshoe/scripts/eventwitter
command=npm start
